# claudecode/rootfs/etc/s6-overlay/s6-rc.d/ttyd/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "ttyd service stopped with exit code: ${1}"

# Give time before restart
sleep 2
